# tools/tg-monitor/deploy/cortex-daily.timer
[Unit]
Description=Run Cortex Daily Digest every morning

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
