# tools/steam-sniper/deploy/steam-sniper-bot.service
[Unit]
Description=Steam Sniper Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/steam-sniper
ExecStart=/root/.local/bin/uv run python main.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/opt/steam-sniper/.env

[Install]
WantedBy=multi-user.target
