# tools/steam-sniper/deploy/steam-sniper-snapshot.timer
[Unit]
Description=Refresh Steam Sniper listings snapshot hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true
Unit=steam-sniper-snapshot.service

[Install]
WantedBy=timers.target
